# picoclaw/picoclaw.conf
[program:picoclaw]
command=/usr/local/bin/picoclaw gateway
autostart=true
autorestart=true
stderr_logfile=/tmp/picoclaw.err.log
stdout_logfile=/tmp/picoclaw.out.log
